# contrib/init/auroracoind.conf
description "Auroracoin Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env AURORACOIND_BIN="/usr/bin/auroracoind"
env AURORACOIND_USER="auroracoin"
env AURORACOIND_GROUP="auroracoin"
env AURORACOIND_PIDDIR="/var/run/auroracoind"
# upstart can't handle variables constructed with other variables
env AURORACOIND_PIDFILE="/var/run/auroracoind/auroracoind.pid"
env AURORACOIND_CONFIGFILE="/etc/auroracoin/auroracoin.conf"
env AURORACOIND_DATADIR="/var/lib/auroracoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # digibyted will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$AURORACOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run auroracoind."
        echo "The setting must appear in $AURORACOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Auroracoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$AURORACOIND_PIDDIR"
    chmod 0755 "$AURORACOIND_PIDDIR"
    chown $AURORACOIND_USER:$AURORACOIND_GROUP "$AURORACOIND_PIDDIR"
    chown $AURORACOIND_USER:$AURORACOIND_GROUP "$AURORACOIND_CONFIGFILE"
    chmod 0660 "$AURORACOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$AURORACOIND_PIDFILE" \
    --chuid $AURORACOIND_USER:$AURORACOIND_GROUP \
    --exec "$AURORACOIND_BIN" \
    -- \
    -pid="$AURORACOIND_PIDFILE" \
    -conf="$AURORACOIND_CONFIGFILE" \
    -datadir="$AURORACOIND_DATADIR" \
    -disablewallet \
    -daemon

